# projects/radmon.conf
[program:radmon]
command=/root/headlessMonitor.py